# .config/autostart/compositor.desktop
[Desktop Entry]
Type=Application
Name=Compositor
Exec=picom -e 1 -i 1 -D 4 --vsync
Terminal=false
