# .config/autostart/locker.desktop
[Desktop Entry]
Type=Application
Name=Locker
Exec=xss-lock -q slock
Terminal=false
