# .config/autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=Unclutter
Exec=unclutter -root
Terminal=false
